# config/helixscreen-update.service
# SPDX-License-Identifier: GPL-3.0-or-later
#
# systemd service — restart HelixScreen after Moonraker update
#
# Triggered by helixscreen-update.path when release_info.json changes.
# Waits 10 seconds for Moonraker to finish extracting all files before
# restarting the service.
#
# After restarting, re-arms the path unit's inotify watch.  Moonraker's
# shutil.rmtree() invalidates watches on the old inode, so we restart
# the path unit to pick up the new file at its new inode.
#
# Workaround for mainsail-crew/mainsail#2444.

[Unit]
Description=Restart HelixScreen after update

[Service]
Type=oneshot

# Skip restart if system booted less than 5 minutes ago.  NTP clock sync
# during early boot can cause PathChanged to fire spuriously (#536).
# A real Moonraker update can't land that soon after power-on.
# Uses /proc/uptime (monotonic, immune to clock jumps).
# Previous 60s-since-helixscreen-start grace was insufficient on systems
# with many MCUs where NTP sync arrives 60-120s+ after boot.
ExecStartPre=/bin/sh -c '\
  UPTIME_S=$$(cut -d. -f1 /proc/uptime); \
  [ "$$UPTIME_S" -lt 300 ] && { echo "Skipping: system booted $${UPTIME_S}s ago (< 300s grace)"; exit 1; } || true'

ExecStartPre=/bin/sleep 10

# Skip restart if helix-screen is handling the update itself (self-update).
# The app writes a sentinel before _exit(0), and the watchdog restarts the
# child internally.  The sentinel is NOT deleted — the path watcher can fire
# multiple times after an atomic directory swap (inotify invalidation +
# delayed systemd rescan), so the sentinel must persist to block all
# spurious triggers.  Cleaned up by the app on next successful startup.
#
# Check two locations:
#   1. @@INSTALL_PARENT@@/.helixscreen/ — preferred, works with PrivateTmp
#   2. /tmp/helixscreen_self_restart — legacy fallback (broken with PrivateTmp)
ExecStartPre=/bin/sh -c '\
  for f in "@@INSTALL_PARENT@@/.helixscreen/self_restart_sentinel" \
           /tmp/helixscreen_self_restart; do \
    if [ -f "$$f" ]; then \
      echo "Skipping: self-update restart handled by watchdog ($$f)"; \
      exit 1; \
    fi; \
  done'

# Refresh systemd units from the newly extracted install dir.
# The template logic lives in a script so Moonraker's extraction
# updates it — avoids the chicken-and-egg where the installer's
# sed pass corrupts @@placeholder@@ patterns embedded in this unit.
# Runs as root (no User= set).
ExecStartPre=@@INSTALL_DIR@@/config/refresh-service-units.sh

ExecStart=/bin/systemctl restart helixscreen
ExecStartPost=/bin/systemctl restart helixscreen-update.path
