# config/helixscreen-update.service
# SPDX-License-Identifier: GPL-3.0-or-later
#
# systemd service — restart HelixScreen after Moonraker update
#
# Triggered by helixscreen-update.path when release_info.json changes.
# Waits 10 seconds for Moonraker to finish extracting all files before
# restarting the service.
#
# After restarting, re-arms the path unit's inotify watch.  Moonraker's
# shutil.rmtree() invalidates watches on the old inode, so we restart
# the path unit to pick up the new file at its new inode.
#
# Workaround for mainsail-crew/mainsail#2444.

[Unit]
Description=Restart HelixScreen after update

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10

# Refresh main service file from the newly extracted install dir.
# Picks up new directives (e.g. StateDirectory) that self-updates
# can't apply due to NoNewPrivileges.  Runs as root (no User= set).
ExecStartPre=/bin/sh -c '\
  SRC="@@INSTALL_DIR@@/config/helixscreen.service"; \
  DEST="/etc/systemd/system/helixscreen.service"; \
  if [ -f "$$SRC" ] && [ -f "$$DEST" ]; then \
    USER=$$(grep "^User=" "$$DEST" | cut -d= -f2); \
    GROUP=$$(grep "^Group=" "$$DEST" | cut -d= -f2); \
    WDIR=$$(grep "^WorkingDirectory=" "$$DEST" | cut -d= -f2); \
    IDIR=$$(dirname "$$(dirname "$$SRC")"); \
    cp "$$SRC" "$$DEST"; \
    sed -i "s|@@HELIX_USER@@|$${USER:-root}|g; s|@@HELIX_GROUP@@|$${GROUP:-root}|g; s|@@INSTALL_DIR@@|$${IDIR}|g" "$$DEST"; \
    systemctl daemon-reload; \
  fi'

ExecStart=/bin/systemctl restart helixscreen
ExecStartPost=/bin/systemctl restart helixscreen-update.path
