# config/helixscreen-update.service
# SPDX-License-Identifier: GPL-3.0-or-later
#
# systemd service — restart HelixScreen after Moonraker update
#
# Triggered by helixscreen-update.path when release_info.json changes.
# Waits 10 seconds for Moonraker to finish extracting all files before
# restarting the service.
#
# After restarting, re-arms the path unit's inotify watch.  Moonraker's
# shutil.rmtree() invalidates watches on the old inode, so we restart
# the path unit to pick up the new file at its new inode.
#
# Workaround for mainsail-crew/mainsail#2444.

[Unit]
Description=Restart HelixScreen after update

[Service]
Type=oneshot

# Skip restart if helixscreen started less than 60s ago.  NTP clock sync
# during early boot can cause PathChanged to fire spuriously (#536).
# A real Moonraker update can't land that fast after boot.
# Uses monotonic timestamps (immune to clock jumps).
ExecStartPre=/bin/sh -c '\
  ACTIVE_MONO=$$(/bin/systemctl show helixscreen.service -p ActiveEnterTimestampMonotonic --value 2>/dev/null); \
  [ "$${ACTIVE_MONO:-0}" = "0" ] && exit 0; \
  UPTIME_S=$$(cut -d. -f1 /proc/uptime); \
  ACTIVE_S=$$((ACTIVE_MONO / 1000000)); \
  ELAPSED=$$((UPTIME_S - ACTIVE_S)); \
  [ "$$ELAPSED" -lt 60 ] && { echo "Skipping: helixscreen started $${ELAPSED}s ago (< 60s grace)"; exit 1; } || true'

ExecStartPre=/bin/sleep 10

# Refresh systemd units from the newly extracted install dir.
# The template logic lives in a script so Moonraker's extraction
# updates it — avoids the chicken-and-egg where the installer's
# sed pass corrupts @@placeholder@@ patterns embedded in this unit.
# Runs as root (no User= set).
ExecStartPre=@@INSTALL_DIR@@/config/refresh-service-units.sh

ExecStart=/bin/systemctl restart helixscreen
ExecStartPost=/bin/systemctl restart helixscreen-update.path
